# debian/mtp-server.conf
description "MTP protocol server"

start on :sys:android-mtp-on or :sys:android-usb-connected
stop on desktop-end or :sys:android-mtp-off

env GOOGLE_LOGTOSTDERR=1

respawn

exec /usr/bin/mtp-server
