# debian/mtp-server.conf
description "MTP protocol server"

start on :sys:android-usb-connected
stop on :sys:android-usb-disconnected

respawn

exec /usr/bin/mtp-server
