# debian/mtp-server.conf
description "MTP protocol server"

start on startup

respawn

exec /usr/bin/mtp-server
